# epox/rootfs-overlay/etc/init.d/calcium-save
#!/sbin/openrc-run

description="Calcium OS Persistence Manager"

depend() {
    need localmount
    use bootmisc
}

start() {
    # Skip if installed system (not live session)
    /usr/bin/calcium-save is-live || { eend 0 "Installed system, nothing to do"; return 0; }
    ebegin "Loading persistent session"
    /usr/bin/calcium-save load || true
    eend 0
}

stop() {
    /usr/bin/calcium-save is-live || { return 0; }
    ebegin "Saving session to boot media"
    /usr/bin/calcium-save save || true
    eend 0
}
