# epox/rootfs-overlay/etc/init.d/tailscale-ssh
#!/sbin/openrc-run

name="tailscale-ssh"
description="Enable Tailscale SSH"

depend() {
    need tailscale
}

start() {
    ebegin "Enabling Tailscale SSH"
    /usr/bin/tailscale set --ssh
    eend $?
}
