# epox/rootfs-overlay/etc/profile.d/local-bin.sh
_local_bin="${HOME}/.local/bin"
[ -d "$_local_bin" ] && PATH="${_local_bin}:${PATH}"
unset _local_bin
